# conf.d/mdns-repeater.conf
[program:mdns-repeater]
user=%(ENV_SYS_USERNAME)s
command=mdns-repeater -f eth0 wg0
autorestart=true
