# conf.d/wg-server.conf
[program:wg-server]
command=/app/run-wg-server.sh
autorestart=true
